# deployments/systemd/change-society-api.service.template
[Unit]
Description=Change Society API (AgentCore)
After=network-online.target change-society-postgres.service
Wants=network-online.target change-society-postgres.service
Requires=change-society-postgres.service

[Service]
Type=simple
WorkingDirectory=@PACK_ROOT@
EnvironmentFile=@PACK_ROOT@/.env
Environment=PYTHONPATH=@PACK_ROOT@/backend/change-society-service/src
Environment=CHANGE_SOCIETY_PACK_ROOT=@PACK_ROOT@
ExecStart=@PYTHON@ @PACK_ROOT@/run.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
