# deployments/systemd/change-society-web.service.template
[Unit]
Description=Change Society Web UI (AgentCore hackathon)
After=network-online.target change-society-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@/hackathon/frontend
EnvironmentFile=@REPO_ROOT@/hackathon/.env
ExecStart=@NPM@ run start
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
