# deployments/systemd/change-society-web.service.template
[Unit]
Description=Change Society Web UI (AgentCore)
After=network-online.target change-society-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@PACK_ROOT@/frontend
EnvironmentFile=@PACK_ROOT@/.env
ExecStart=@NPM@ run start
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
